# L1.V27/deploy/systemd/mywebapp.socket
[Unit]
Description=mywebapp listening socket (variant 27, localhost only)
PartOf=mywebapp.service

[Socket]
ListenStream=127.0.0.1:3000
NoDelay=true
SocketUser=app
SocketGroup=app
SocketMode=0660

[Install]
WantedBy=sockets.target
